# {{cookiecutter.project_name}}/deploy/supervisor_django.conf
[program:django]
command=/app/venv/bin/uwsgi
    --chdir /app
    --env DJANGO_SETTINGS_MODULE=settings.production
    --env PYTHONHASHSEED=random
    --env LANG=en_US.UTF-8
    --socket /app/uwsgi.sock
    --chmod-socket
    --module {{ cookiecutter.project_name }}.wsgi:application
    --master
    --processes 10
    --no-orphans
    --vacuum
    --log-4xx
    --log-5xx
    --single-interpreter
    --log-zero
    --log-slow 1000
    --virtualenv /app/venv
    --logdate
environment=PYTHONHASHSEED="random"
directory=/app
numprocs=1
stdout_logfile=/var/log/django_app.log
autostart=true
autorestart=true
startsecs=10
startretries=5
redirect_stderr=true
stopsignal=QUIT
stopwaitsecs=60

[program:nginx]
command = /usr/sbin/nginx
